# etc/services.d/logstash/run
#!/usr/bin/with-contenv /bin/ash
exec s6-setuidgid logstash /opt/containerpilot/containerpilot /usr/share/logstash/bin/logstash agent --allow-env -f /usr/share/logstash/config/logstash.conf